# waybar-notes-sync.timer
[Unit]
Description=Waybar Notes - Google Tasks Sync Timer
Requires=waybar-notes-sync.service

[Timer]
# Sync every 10 minutes
OnBootSec=1min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
